# monitor-server.cron
# Runs every three hours to poll all nodes and PCUs, as well as take some
# automated actions for debug nodes.

01 * * * * root /usr/share/monitor/automate-default.sh > /usr/share/monitor/monitor.log 2>&1
#30 * * * * root /etc/plc.d/monitor sync >> /var/log/monitorsync.log 2>&1

