# nagios/monitor-nagios.cron
# run daily to regenerate the nagios configuration files
0 0 * * * root /usr/share/monitor/nagios/plc_hosts_to_nagios.py > /etc/nagios/objects/plcnodes.cfg
5 0 * * * root /usr/share/monitor/nagios/plc_users_to_nagios.py > /etc/nagios/objects/plcusers.cfg
10 0 * * * root /sbin/service nagios restart
