# rt3/syncrtusers.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/share/monitor/rt3
HOME=/usr/share/monitor
MAILTO=soltesz@cs.princeton.edu

01 * * * * root rtcron.sh

